# build/docker/httpWatchDemo/conf.d/client.ini
[program:juno-agent]
directory=/home/www/server/juno-agent
command=/home/www/server/juno-agent/bin/juno-agent --config=./config/config-local.toml


[program:client]
directory=/home/www/server/client
command=/home/www/server/client/bin/httpWatchDemo --config="http://127.0.0.1:50010/api/v1/agent/config?name=douyu-edogawa-konan&env=live&target=config-live.toml&port=8013" --watch
