# build/docker/fileWatchDemo/conf.d/client.ini
[program:juno-agent]
directory=/home/www/server/juno-agent
command=/home/www/server/juno-agent/bin/juno-agent --config=./config/config-local.toml


[program:client]
directory=/home/www/server/client
command=/home/www/server/client/bin/fileWatchDemo --config=/tmp/www/server/douyu-edogawa-konan/config/config-dev.toml --watch=true
